# systemd/vigild.service
[Unit]
Description=Vigil Baseline — Lightweight File Integrity Monitor
Documentation=https://github.com/lousclues-labs/vigil
After=local-fs.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/vigild
# /usr/bin/kill (not /bin/kill): on usrmerge systems both resolve to
# the same inode, but minimal Fedora container images can have a
# broken /bin symlink at systemd-analyze verify time, which then
# rejects the unit with 'Command /bin/kill is not executable'.
ExecReload=/usr/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vigil /run/vigil
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
MemoryMax=50M
MemoryHigh=30M

# Required capabilities
#
# CAP_SYS_ADMIN       -- fanotify_init / fanotify_mark
# CAP_DAC_READ_SEARCH -- bypass DAC on read across uid/mount boundaries
#
# Provided two ways for defence in depth:
#   1. File caps on /usr/bin/vigild via setcap in the postinst
#      scriptlet. Lets manual `vigild` invocations (testing,
#      troubleshooting) work without sudo for fanotify-only ops.
#   2. AmbientCapabilities= here. systemd grants these to the exec'd
#      process directly. If the postinst's setcap fails or is
#      reverted (some minimal container images break the binary's
#      --help after setcap because file caps trigger AT_SECURE=1
#      secure-execution mode in the dynamic linker), the systemd-
#      launched vigild still gets the caps it needs.
#
# CapabilityBoundingSet= caps both paths -- the process can never
# acquire any capability outside this set, regardless of source.
# vigild spawns no children, so the ambient set leaking is moot.
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH

# Watchdog
WatchdogSec=120

# Allow generous startup time for initial baseline scan
TimeoutStartSec=300

# PID file
PIDFile=/run/vigil/vigild.pid
RuntimeDirectory=vigil
RuntimeDirectoryMode=0750
StateDirectory=vigil
StateDirectoryMode=0750
LogsDirectory=vigil
LogsDirectoryMode=0750

[Install]
WantedBy=multi-user.target
