# systemd/vigild.service
[Unit]
Description=VigilBaseline — Lightweight File Integrity Monitor
Documentation=man:vigild(8)
After=local-fs.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/vigild
ExecReload=/bin/kill -HUP $MAINPID
Restart=on-failure
RestartSec=5

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/vigil /run/vigil
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes

# Resource limits
MemoryMax=50M
MemoryHigh=30M

# Required capabilities
AmbientCapabilities=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH
CapabilityBoundingSet=CAP_SYS_ADMIN CAP_DAC_READ_SEARCH

# Watchdog
WatchdogSec=120

# Allow generous startup time for initial baseline scan
TimeoutStartSec=300

# PID file
PIDFile=/run/vigil/vigild.pid
RuntimeDirectory=vigil
RuntimeDirectoryMode=0750
StateDirectory=vigil
StateDirectoryMode=0750
LogsDirectory=vigil
LogsDirectoryMode=0750

[Install]
WantedBy=multi-user.target
