# systemd/vigil-scan.timer
[Unit]
Description=Vigil Baseline scheduled integrity scan
Documentation=https://github.com/lousclues-labs/vigil

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=900
Persistent=true

[Install]
WantedBy=timers.target
